# rootfs/etc/services.d/km-apiserver/run
#!/usr/bin/with-contenv sh
exec 2>&1
${K8SM_HOME:-/opt/k8sm}/km apiserver \
  --insecure-bind-address=0.0.0.0 \
  --insecure-port=${K8SM_PORT:-8888} \
  --etcd-servers=http://${ETCD_HOST:-etcd.service.consul}:${ETCD_PORT:-4001} \
  --service-cluster-ip-range=${K8SM_PORTAL:-10.254.0.0/24} \
  --cloud-provider=mesos \
  --cloud-config=${K8SM_HOME:-/opt/k8sm}/mesos-cloud.conf \
  --v=1
