# rootfs/etc/services.d/km-controller/run
#!/usr/bin/with-contenv sh
exec 2>&1
km controller-manager \
  --master=${K8SM_HOST}:${K8SM_PORT} \
  --cloud-provider="mesos" \
  --cloud-config=${K8SM_HOME}/mesos-cloud.conf
