# rootfs/etc/services.d/km-controller/run
#!/usr/bin/with-contenv sh
exec 2>&1
km controller-manager \
  --master=${K8SM_IP}:${K8SM_PORT} \
  --cloud_config=${K8SM_CONFIG}